# slack-bot.service
[Unit]
Description=Multi-Function Slack Bot
After=network.target

[Service]
Type=simple
User=molt
WorkingDirectory=/home/molt/Slack Bot
ExecStart=/home/molt/Slack Bot/.venv/bin/python /home/molt/Slack Bot/main.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
