# slack-bot.service
[Unit]
Description=Multi-Bot Slack Service
After=network.target

[Service]
Type=simple
User=molt
WorkingDirectory=/home/molt/Slack Bot
ExecStart=/home/molt/Slack\x20Bot/.venv/bin/python /home/molt/Slack\x20Bot/run_all.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
